# rootfs/etc/services.d/fhem/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: FHEM
# Starts server to serve the app
# ==============================================================================
bashio::log.info "Run FHEM..."

cd /opt/fhem \
    || bashio::exit.nok "Could not change into app directory"

exec perl fhem.pl fhem.cfg